# services/eth-bootnode-healthcheck.service
[Unit]
Description=eth-bootnode-healthcheck
After=network.target auditd.service
Wants=network.target

[Service]
ExecStart=/usr/bin/eth-bootnode-healthcheck-linux-x64
Restart=always
RestartSec=5s
Environment="MAX_BLOCK_DIFFERENCE=3"

[Install]
WantedBy=default.target 
RequiredBy=swarm.service
Alias=eth-bootnode-healthcheck
